# deploy/wealthtrack.service
[Unit]
Description=WealthTrack API
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=filla
WorkingDirectory=/home/filla/dev/wealthtrack/backend
Environment=PATH=/home/filla/.local/bin:/home/filla/dev/wealthtrack/.venv/bin
ExecStart=/home/filla/dev/wealthtrack/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
